# dispatcharr/rootfs/etc/services.d/daphne/run
#!/usr/bin/with-contenv bashio
cd /app
export PYTHONPATH="/app"
export $(grep -v '^#' /app/.env | xargs)
exec daphne -b 0.0.0.0 -p 8001 dispatcharr.asgi:application
