# dispatcharr/rootfs/etc/services.d/daphne/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Daphne (WebSocket server) on port 8001"

cd /app || exit 1

export PYTHONPATH="/app"
export PATH="/app/env/bin:$PATH"

# Load environment safely
set -a
. /app/.env
set +a

exec /app/env/bin/daphne \
    -b 0.0.0.0 \
    -p 8001 \
    dispatcharr.asgi:application
