# dispatcharr/rootfs/etc/services.d/celerybeat/run
#!/usr/bin/with-contenv bashio
cd /app
export PYTHONPATH="/app"
export PATH="/app/env/bin:$PATH"

# Load your Django/Database settings
if [ -f /app/.env ]; then
    export $(grep -v '^#' /app/.env | xargs)
fi
exec celery -A dispatcharr beat -l info
