# dispatcharr/rootfs/etc/services.d/redis/run
cd /app
export PYTHONPATH="/app"
export $(grep -v '^#' /app/.env | xargs)
#!/usr/bin/with-contenv bashio
exec redis-server --daemonize no
