# dispatcharr/rootfs/etc/services.d/redis/run
cd /app
export PYTHONPATH="/app"
export PATH="/app/env/bin:$PATH"

# Load environment safely
set -a
. /app/.env
set +a

#!/usr/bin/env bashio
exec redis-server /etc/redis.conf
