# dispatcharr/rootfs/etc/services.d/celery/run
#!/usr/bin/with-contenv bashio
export $(grep -v '^#' /app/.env | xargs)
cd /app
exec celery -A dispatcharr worker -l info
