# dispatcharr/rootfs/etc/services.d/postgres/run
cd /app
export PYTHONPATH="/app"
export PATH="/app/env/bin:$PATH"

# Load environment safely
set -a
. /app/.env
set +a

#!/usr/bin/env bashio
exec su-exec postgres postgres -D /data/db -c unix_socket_directories='/run/postgresql'
