# dispatcharr/rootfs/etc/services.d/postgres/run
cd /app
export PYTHONPATH="/app"
export PATH="/app/env/bin:$PATH"

# Load your Django/Database settings
if [ -f /app/.env ]; then
    export $(grep -v '^#' /app/.env | xargs)
fi
#!/usr/bin/env bashio
exec su-exec postgres postgres -D /data/db -c unix_socket_directories='/run/postgresql'
