# dispatcharr/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
exec su-exec postgres postgres -D /data/db
