# dispatcharr/rootfs/etc/services.d/postgres/run
cd /app
export PYTHONPATH="/app"
export $(grep -v '^#' /app/.env | xargs)
#!/usr/bin/with-contenv bashio
exec su-exec postgres postgres -D /data/db
