# dispatcharr/rootfs/etc/services.d/gunicorn/run
#!/usr/bin/with-contenv bashio

# Fetch the port from the user's config
WEB_PORT=$(bashio::config 'web_port')

# Fetch the key from the user's config
export SECRET_KEY="8Ebz9PHj1n95CHnj9"

bashio::log.info "Starting Dispatcharr on host port: ${WEB_PORT}"

# Bind Gunicorn to that specific port
exec /app/env/bin/gunicorn dispatcharr.wsgi:application --bind 0.0.0.0:${WEB_PORT}
