# dispatcharr/rootfs/etc/services.d/gunicorn/run
#!/usr/bin/with-contenv bashio

# Load environment (critical)
set -a
. /app/.env
set +a

#!/usr/bin/with-contenv bashio
# Wait for postgres socket
while [ ! -S /run/postgresql/.s.PGSQL.5432 ]; do
  sleep 1
done

cd /app
exec /app/env/bin/gunicorn --timeout 600 dispatcharr.wsgi:application
