# dispatcharr/rootfs/etc/services.d/gunicorn/run
#!/usr/bin/with-contenv bashio

cd /app

# Load environment (critical)
set -a
. /app/.env
set +a

# Wait for postgres socket
while [ ! -S /run/postgresql/.s.PGSQL.5432 ]; do
  sleep 1
done

exec /app/env/bin/gunicorn \
  --timeout 1200 \
  --workers 2 \
  --worker-class gthread \
  --threads 4 \
  --bind 0.0.0.0:8000 \
  dispatcharr.wsgi:application
