# dispatcharr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Nginx on port 9191"

exec nginx -g "daemon off;"
